# recipes/laundry.service
[Unit]
Description=Laundry service

[Service]
ExecStart=/home/ubuntu/src/recipes/start.sh

[Install]
WantedBy=multi-user.target
